# KubeArmor/packaging/kubearmor.service
[Unit]
Description=KubeArmor

[Service]
User=root
KillMode=process
WorkingDirectory=/opt/kubearmor/
ExecStart=/opt/kubearmor/kubearmor

[Install]
WantedBy=multi-user.target
